# ubuntu-14.04/etc/cron.d/chros73
HOMEBIN="/home/chros73/bin"

# run backup script every night at 01:17 , dismiss any errors
17 1 * * * 	chros73	"$HOMEBIN/doBackup.sh" > /dev/null 2>&1
# run listing all orphaned data in every 10 days at 04:48
48 4 */10 * *	chros73	"$HOMEBIN/listOrphans.sh"

